# supervisord.conf
; 三个进程都交给 supervisord 托管，目的是自愈：
; warp-svc 挂掉时如果没人拉起，容器还活着但代理已经静默失效——
; 那是最难排查的一类故障。
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid
user=root

; 没有这三段，容器里的 supervisorctl 会报
; "unix:///var/run/supervisor.sock no such file"，排障手段直接失效。
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; warp-svc 需要系统 dbus。不跑 dbus 功能上也能用，但 power_notifier 会
; 每 3 秒打一条「dbus connection failed」告警且永不停止（实测约 2.9 万行/天），
; 足以把 docker logs 冲垮。装上 dbus 是比过滤日志更正确的解法。
[program:dbus]
command=/bin/bash -c "mkdir -p /run/dbus && rm -f /run/dbus/pid && exec dbus-daemon --system --nofork --nopidfile"
priority=5
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

; 经包装器启动以过滤 DEBUG 刷屏（warp-svc 不认 RUST_LOG，详见脚本内注释）。
; stopasgroup/killasgroup 是必需的：包装器里是管道，信号要送到整个进程组，
; 否则 warp-svc 收不到 SIGTERM，停容器时会被硬杀。
[program:warp-svc]
command=/opt/run-warp-svc.sh
priority=10
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

; 注册 / 设协议 / 连接 / 装策略路由，之后进入看护循环常驻。
; autorestart=true 是为了 warp-svc 重启后能重建连接与路由。
[program:warp-agent]
command=/opt/warp-agent.sh
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:xray]
command=/usr/local/bin/xray run -c /etc/xray/config.json
priority=30
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
